# configs/auto-reboot.service
[Unit]
Description=Reboot if a kernel update is pending
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'if [ -f /var/run/reboot-required ]; then /sbin/shutdown -r +1 "Rebooting: kernel update pending"; fi'
